# ansible/templates/kiosk.service.j2
[Unit]
Description=Kiosk Browser
After=multi-user.target seatd.service ems-openhours.service
Requires=seatd.service ems-openhours.service

[Service]
Type=simple
User={{ kiosk_user }}
Group={{ kiosk_user }}

# Ensure XDG_RUNTIME_DIR exists (logind won't create it for system users)
RuntimeDirectory=kiosk
Environment=XDG_RUNTIME_DIR=/run/kiosk
Environment=LIBSEAT_BACKEND=seatd
Environment=WLR_LIBINPUT_NO_DEVICES=1
Environment=WLR_NO_HARDWARE_CURSORS=1

# Kiosk configuration
Environment=KIOSK_URL={{ kiosk_url }}

# cage runs Firefox in kiosk mode on a dedicated TTY
# wlr-randr applies rotation after cage has started
ExecStart=/bin/sh -c '(sleep 3 && wlr-randr --output HDMI-A-1 --transform {{ display_rotation }}) & exec /usr/bin/cage -- /usr/bin/firefox-esr --kiosk {{ kiosk_url }}'

# Restart on crash — 5s delay to avoid rapid restart loops
Restart=on-failure
RestartSec=5

# Run on TTY7 (typical for graphical sessions)
TTYPath=/dev/tty7
StandardInput=tty
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
